# scripts/systemd/ibkr-gateway.service
[Unit]
Description=IBKR Gateway (IBC)
After=network.target

[Service]
Type=forking
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/pearlalgo-dev-ai-agents
Environment="DISPLAY=:99"
Environment="PATH=/home/pearlalgo/pearlalgo-dev-ai-agents/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/usr/bin/Xvfb :99 -screen 0 1024x768x24 &
ExecStart=/home/pearlalgo/pearlalgo-dev-ai-agents/scripts/gateway/gateway.sh start
ExecStop=/home/pearlalgo/pearlalgo-dev-ai-agents/scripts/gateway/gateway.sh stop
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
