# scripts/systemd/pearlalgo-api-tv-paper.service
[Unit]
Description=PearlAlgo API Server - Tradovate Paper (port 8001)
After=network.target pearlalgo-agent-tv-paper.service
Wants=pearlalgo-agent-tv-paper.service

[Service]
Type=simple
User=pearl
Group=pearl
WorkingDirectory=/home/pearl/PearlAlgoProject
Environment="PATH=/home/pearl/PearlAlgoProject/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pearl/PearlAlgoProject/.env
EnvironmentFile=-/home/pearl/.config/pearlalgo/secrets.env
ExecStart=/home/pearl/PearlAlgoProject/.venv/bin/python scripts/pearlalgo_web_app/api_server.py --data-dir data/tradovate/paper --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
