# scripts/systemd/pearlalgo-webapp.service
[Unit]
Description=PearlAlgo Web App (Next.js)
After=network.target pearlalgo-api.service
Wants=pearlalgo-api.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/projects/pearl-algo/apps/pearl-algo-app
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="PORT=3001"
Environment="HOSTNAME=0.0.0.0"
EnvironmentFile=/home/pearlalgo/projects/pearl-algo/.env
EnvironmentFile=-/home/pearlalgo/.config/pearlalgo/secrets.env
EnvironmentFile=-/home/pearlalgo/projects/pearl-algo/apps/pearl-algo-app/.env.local
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
