# scripts/systemd/pearlalgo-telegram.service
[Unit]
Description=PearlAlgo Telegram Handler
After=network.target pearlalgo-agent.service
Wants=pearlalgo-agent.service

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/PearlAlgoWorkspace
Environment="PATH=/home/pearlalgo/PearlAlgoWorkspace/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PEARLALGO_API_URL=http://localhost:8001"
EnvironmentFile=/home/pearlalgo/PearlAlgoWorkspace/.env
EnvironmentFile=-/home/pearlalgo/.config/pearlalgo/secrets.env
ExecStart=/home/pearlalgo/PearlAlgoWorkspace/.venv/bin/python -m pearlalgo.telegram.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
