# scripts/systemd/pearlalgo-api.service
[Unit]
Description=PearlAlgo API Server - Tradovate Paper (port 8001)
After=network.target pearlalgo-agent.service
Wants=pearlalgo-agent.service

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/PearlAlgoWorkspace
Environment="PATH=/home/pearlalgo/PearlAlgoWorkspace/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pearlalgo/PearlAlgoWorkspace/.env
EnvironmentFile=-/home/pearlalgo/.config/pearlalgo/secrets.env
ExecStart=/home/pearlalgo/PearlAlgoWorkspace/.venv/bin/python scripts/pearlalgo_web_app/api_server.py --data-dir data/tradovate/paper --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
