# scripts/systemd/pearlalgo-monitor.timer
[Unit]
Description=PearlAlgo Health Monitor Timer
Requires=pearlalgo-monitor.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
