# scripts/systemd/pearlalgo-agent.service
[Unit]
Description=PearlAlgo Market Agent - IBKR Virtual (NQ)
After=network.target ibkr-gateway.service
Wants=ibkr-gateway.service

[Service]
Type=simple
User=pearl
Group=pearl
WorkingDirectory=/home/pearl/PearlAlgoProject
Environment="PATH=/home/pearl/PearlAlgoProject/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pearl/PearlAlgoProject/.env
ExecStart=/home/pearl/PearlAlgoProject/.venv/bin/python -m pearlalgo.market_agent.service --market NQ
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
