# scripts/systemd/pearlalgo-telegram.service
[Unit]
Description=PearlAlgo Telegram Handler
After=network.target pearlalgo-agent.service
Wants=pearlalgo-agent.service

[Service]
Type=simple
User=pearl
Group=pearl
WorkingDirectory=/home/pearl/PearlAlgoProject
Environment="PATH=/home/pearl/PearlAlgoProject/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pearl/PearlAlgoProject/.env
ExecStart=/home/pearl/PearlAlgoProject/.venv/bin/python -m pearlalgo.market_agent.telegram_command_handler
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
