# scripts/systemd/pearlalgo-agent.service
[Unit]
Description=PearlAlgo Market Agent - Tradovate Paper
After=network.target ibkr-gateway.service
Wants=ibkr-gateway.service

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/projects/pearl-algo
Environment="PATH=/home/pearlalgo/projects/pearl-algo/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/pearlalgo/projects/pearl-algo/src"
EnvironmentFile=/home/pearlalgo/projects/pearl-algo/.env
EnvironmentFile=-/home/pearlalgo/.config/pearlalgo/secrets.env
ExecStart=/home/pearlalgo/projects/pearl-algo/.venv/bin/python -m pearlalgo.market_agent.main --config config/live/tradovate_paper.yaml --data-dir data/agent_state/MNQ
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
