# scripts/automated_trading.service
[Unit]
Description=PearlAlgo Automated Trading Agent (IBKR Paper Trading)
After=network.target ibgateway.service
Wants=ibgateway.service

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/pearlalgo-dev-ai-agents
Environment="PATH=/home/pearlalgo/pearlalgo-dev-ai-agents/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pearlalgo/pearlalgo-dev-ai-agents/.venv/bin/python scripts/automated_trading.py --symbols ES NQ GC --strategy sr --interval 300 --log-file logs/automated_trading.log
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=2G
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

