# scripts/ibgateway.service
# Legacy standalone launcher (without IBC). Prefer ibgateway-ibc.service.example for auto-login.
# Copy to /etc/systemd/system/ibgateway.service (requires sudo) then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable ibgateway.service
#   sudo systemctl start ibgateway.service
#
# Notes:
# - Expects IB Gateway installed under /home/pearlalgo/Jts/ibgateway/1041
# - Requires xvfb-run installed (sudo apt-get install xvfb)
# - API settings and auto-login must already be saved in ~/Jts (via GUI or IBC)

[Unit]
Description=IB Gateway (headless) for PearlAlgo
After=network.target

[Service]
Type=simple
User=pearlalgo
WorkingDirectory=/home/pearlalgo/Jts/ibgateway/1041
ExecStart=/usr/bin/xvfb-run -a /home/pearlalgo/Jts/ibgateway/1041/ibgateway
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
