# scripts/pearlalgo_signal_service.service
[Unit]
Description=PearlAlgo 24/7 Signal Generation Service (Polygon)
After=network.target
Wants=network.target

[Service]
Type=simple
User=pearlalgo
Group=pearlalgo
WorkingDirectory=/home/pearlalgo/pearlalgo-dev-ai-agents
Environment="PATH=/home/pearlalgo/pearlalgo-dev-ai-agents/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Main service command
ExecStart=/home/pearlalgo/pearlalgo-dev-ai-agents/.venv/bin/python \
    /home/pearlalgo/pearlalgo-dev-ai-agents/scripts/signal_generation_service.py \
    --symbols ES NQ MES MNQ \
    --strategy sr \
    --interval 300 \
    --log-file logs/signal_generation.log \
    --max-retries 3 \
    --retry-backoff 60.0

# Restart policy
Restart=always
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryMax=2G
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pearlalgo-signals

[Install]
WantedBy=multi-user.target

